# nzbget.conf
[program:nzbget]
user = nobody
startsecs = 0
autorestart = False
command = /home/nobody/start.sh