# supervisor/nzbget.conf
[program:start-script]
autorestart = false
startsecs = 0
user = root
command = /root/start.sh

[program:nzbget]
user = nobody
command = /home/nobody/start.sh
umask = 000
